# online-questions.service
[Unit]
Description=Online Questions
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=BASE_URL=http://test.mattvenn.net
Environment=TEACHER_PASSWORD=changeme
Environment=SECRET_KEY=change-this-to-a-long-random-string
ExecStart=__INSTALL_DIR__/venv/bin/gunicorn --workers 1 --bind 127.0.0.1:5001 server:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
